# inkypi-weather/package/InkyPi/install/privileged/inkypi-privileged.socket
[Unit]
Description=InkyPi fixed privileged action socket

[Socket]
ListenStream=/run/inkypi-privileged.sock
SocketUser=root
SocketGroup=inkypi
SocketMode=0660
RemoveOnStop=true

[Install]
WantedBy=sockets.target
