# inkypi-weather/package/InkyPi/install/inkypi.service
[Unit]
Description=InkyPi App
After=network-online.target inkypi-privileged.socket
Wants=network-online.target
Requires=inkypi-privileged.socket

[Service]
User=inkypi
Group=inkypi
RuntimeDirectory=inkypi
WorkingDirectory=/run/inkypi
EnvironmentFile=-/etc/inkypi/inkypi.env
Environment="PYTHONDONTWRITEBYTECODE=1" "INKYPI_CACHE_DIR=/var/cache/inkypi" "INKYPI_DATA_DIR=/var/lib/inkypi/data" "LG_WD=/run/inkypi"
ExecStart=/usr/local/bin/inkypi run
Restart=on-failure
RestartSec=60
KillSignal=SIGINT
TimeoutStopSec=240
StandardOutput=journal
StandardError=journal
CPUQuota=40%
MemoryMax=200M
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/inkypi /var/cache/inkypi /run/inkypi /etc/inkypi
UMask=0077

[Install]
WantedBy=multi-user.target
